# systemd/mavlink-companion-log-service.service
[Unit]
Description=MAVLink Companion Service
Documentation=https://github.com/Angad7600123/Mavlink-companion-log-service
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=mcls
Group=mcls
ExecStart=/usr/local/bin/mcls /etc/mcls/config.toml
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/mcls
ConfigurationDirectory=mcls
ConfigurationDirectoryMode=0755
StateDirectory=mcls
StateDirectoryMode=0750
AmbientCapabilities=
CapabilityBoundingSet=
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
